# supervisorctl.conf
[program:vehicle_genius_admin]
command=bash start.sh
directory=/home/dimo/vehicle-genius-admin
