# supervisorctl.conf
[program:vehicle_genius_admin]
command=sudo gunicorn vehicle_genius_admin.wsgi -b 127.0.0.1:8001 -w 3
directory=/home/dimo/vehicle-genius-admin
